# amdgpu-overclock.service
[Unit]
Description=amdgpu-overclock

[Service]
Type=oneshot
ExecStart=/usr/bin/amdgpu-overclock

[Install]
WantedBy=multi-user.target
